# scripts/cloudflare/systemd/cloudflare-allowlist-update.timer
[Unit]
Description=Run Cloudflare origin allowlist update daily

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target

